# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
priority=1

[program:backend]
directory=/app/backend
command=npm start
autorestart=true
priority=2
